# services/backend.service.template
[Unit]\nDescription=Malawi Chatbot Backend\nAfter=network.target\n\n[Service]\nType=simple\nUser=<your_username>\nWorkingDirectory=/path/to/chatbot-mukundi\nEnvironment=PYTHONPATH=/path/to/chatbot-mukundi\nEnvironment=PATH=/path/to/chatbot-mukundi/.venv/bin:/usr/local/bin:/usr/bin\nEnvironment=VIRTUAL_ENV=/path/to/chatbot-mukundi/.venv\nExecStart=/bin/bash -c 'cd /path/to/chatbot-mukundi && source .venv/bin/activate && python -m uvicorn app.main:app --host 0.0.0.0 --port 8001'\nRestart=always\nRestartSec=5\n\nStandardOutput=append:/var/log/malawi-backend.log\nStandardError=append:/var/log/malawi-backend.error.log\n\n[Install]\nWantedBy=multi-user.target
